# devspaces/docker/assets/services/postgres/run
#!/bin/bash

set -e
exec 2>&1

exec chpst -U postgres:postgres pg_ctlcluster $PG_VERSION main start --foreground
